# infra/systemd/wasm-pdf-backup.timer
[Unit]
Description=Daily encrypted PostgreSQL backup for WASM PDF Editor

[Timer]
OnCalendar=*-*-* 02:15:00 UTC
RandomizedDelaySec=30m
Persistent=true
Unit=wasm-pdf-backup.service

[Install]
WantedBy=timers.target
