# infra/systemd/wasm-pdf-backup.service
[Unit]
Description=Encrypted PostgreSQL backup for WASM PDF Editor
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
User=deploy
Group=deploy
UMask=0077
WorkingDirectory=/opt/wasm-pdf-editor/current
EnvironmentFile=/etc/wasm-pdf-editor/deployment.env
EnvironmentFile=/etc/wasm-pdf-editor/backup.env
ExecStart=/opt/wasm-pdf-editor/current/scripts/backup-postgres.sh
TimeoutStartSec=2h
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadWritePaths=/var/backups/wasm-pdf-editor
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
